# templates/pull-routing.timer.tmpl
[Unit]
Description=Refresh SIP routing from PostgREST view into Kamailio mtree

[Timer]
OnBootSec=30s
OnUnitActiveSec=__PULL_INTERVAL__

[Install]
WantedBy=timers.target
