# templates/pull-routing.service.tmpl
[Unit]
Description=Export SIP routing from PostgREST view to Kamailio mtree
After=network-online.target kamailio.service
Wants=network-online.target

[Service]
Type=oneshot
User={{KAM_USER}}
Group={{KAM_GROUP}}
EnvironmentFile=-{{REPO_DIR}}/.env
ExecStart={{REPO_DIR}}/bin/pull-routing.sh
# Don't fail systemd if kamcmd can’t reach Kamailio (e.g. first boot)
SuccessExitStatus=0 1

[Install]
WantedBy=multi-user.target
