# templates/pull-routing.timer.tmpl
[Unit]
Description=Periodic PostgREST → Kamailio mtree sync

[Timer]
OnBootSec=30s
OnUnitActiveSec={{PULL_MINUTES}}min
AccuracySec=10s
Unit=pull-routing.service

[Install]
WantedBy=timers.target
